# bukhari-sales.service
[Unit]
Description=Bukhari Sales — AI Sales Bot (Instagram + Telegram)
After=network.target

[Service]
Type=simple
User=qutb
Group=qutb
WorkingDirectory=/home/qutb/bukhari-sales/backend
ExecStart=/home/qutb/bukhari-sales/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PATH=/home/qutb/bukhari-sales/backend/.venv/bin:/usr/bin:/bin
EnvironmentFile=/home/qutb/bukhari-sales/backend/.env

[Install]
WantedBy=multi-user.target
